# tasker-on-dotcloud/taskernodejs/supervisord.conf
[program:node]
command = node server.js
directory = /home/dotcloud/current
environment = VERBOSE=FALSE